# luci-app-netspeedtest/root/etc/uci-defaults/zzz_luci-app-netspeedtest
#!/bin/sh

chmod +x /usr/bin/netspeedtest.sh /usr/bin/speedtest /usr/libexec/rpcd/luci.netspeedtest

rm -f /tmp/luci-indexcache
exit 0
