# luci-app-netspeedtest/root/etc/uci-defaults/zzz_luci-app-netspeedtest
#!/bin/sh

chmod +x /etc/init.d/netspeedtest /usr/bin/netspeedtest.sh /usr/bin/speedtest /usr/libexec/rpcd/luci.netspeedtest
/etc/init.d/rpcd restart
rm -f /tmp/luci-indexcache
exit 0
